# examples/seglake.service
[Unit]
Description=Seglake
After=network.target

[Service]
User=seglake
Group=seglake
WorkingDirectory=/var/lib/seglake
EnvironmentFile=/etc/seglake/secrets.env
Environment=SEGLAKE_ADDR=127.0.0.1:9000
Environment=SEGLAKE_DATA_DIR=/var/lib/seglake
Environment=SEGLAKE_REGION=us-east-1
Environment=SEGLAKE_PUBLIC_BUCKETS=public
Environment=SEGLAKE_TRUSTED_PROXIES=127.0.0.1/32,::1/128
Environment=SEGLAKE_VIRTUAL_HOSTED=false
# Bool flags must use --flag=false (no space) when disabling.
# Each continued line requires a trailing '\' for systemd.
# Prefer env vars over flags to avoid secrets in process listings.
ExecStart=/opt/seglake/build/seglake \
  -mode server \
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
